# dots/config/systemd/user/mon-ctrl-time.service
[Unit]
Description=Adjust monitor brightness based on time of day

[Service]
Type=oneshot
ExecStart=%h/.local/bin/mon-ctrl time
